# deployer-script/deployer.path
[Path]
PathChanged=/var/www/html/wp-content/uploads/deploy.txt

[Unit]
Wants=deployer.service

